# websiteapp/tools/Docker/websiteguide.ini
[supervisord]
nodaemon=true
[program:websiteguide]
command = bash /WebsiteGuide/websiteapp/tools/start-api.sh
autostart = true
stdout_logfile = /WebsiteGuide/websiteapp/logs/api.log
redirect_stderr = true