# websiteapp/tools/WebsiteGuide.ini
[supervisord]
nodaemon=true

[program:spug-api]
command = bash /WebsiteGuide/websiteapp/tools/start-api.sh
autostart = true
stdout_logfile = /WebsiteGuide/api.log
redirect_stderr = true